# docker/supervisor/kdc.conf
[program:kdc]
command=/usr/sbin/krb5kdc -n
autostart=true
autorestart=true
